# pasta-user.service
[Service]
ExecStart=/home/$USER/pasta
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
